# crontab
23 07 * * * cd / && python rfm_scraper.py >> output.log 2>&1 &
02 15 * * * cd / && python rfm_scraper.py >> output.log 2>&1 &